# supervisord-tomcat.conf
[program:tomcat]
command=/opt/apache-tomcat-8.0.47/bin/catalina.sh run
numprocs=1
autostart=true
autorestart=false
user=root
