# supervisord-exist.conf
[program:existdb]
command=/opt/exist/bin/startup.sh
numprocs=1
autostart=true
autorestart=true
user=root
